# deploy/candle-logger.service
[Unit]
Description=OPT_Trader EOD 1-min candle logger (DXFeed, same feed as fills)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# ── Set these to MATCH this box's trading unit (copy from optionsbot.service) ──
User=ubuntu
WorkingDirectory=/home/ubuntu/options-trader
# The logger needs the SAME environment as the bot: Tastytrade creds (for the
# DXLink session) AND OT_INSTRUMENT (so it defaults to this box's symbol). Point
# this at the same source your bot uses. If your bot keeps its env INLINE in its
# systemd unit rather than a file, copy those Environment= lines in here instead.
EnvironmentFile=/home/ubuntu/options-trader/.env
# No --symbols / --out needed: defaults to this box's OT_INSTRUMENT, written to
# <repo>/data/OHLC/<YYYY-MM-DD>/<SYMBOL>.csv. On an SPX box, append the map:
#   ... -m data.candle_logger --symbol-map SPX=SPX
ExecStart=/home/ubuntu/options-trader/venv/bin/python -m data.candle_logger

# Output: <WorkingDirectory>/data/OHLC/<YYYY-MM-DD>/<SYMBOL>.csv
